# test/mkosi/trees/test/etc/systemd/system/sshcad.service.d/10-test-user.conf
[Service]
ExecStartPre=/bin/bash -c 'echo -e "testpass\\ntestpass" | /usr/local/bin/sshcad add-user testuser'
