# test/mkosi/trees/test/etc/systemd/system/sshcad.service.d/10-test-user.conf
[Service]
ExecStartPre=/bin/bash -c 'echo -e "testuser\\ntestpass\\n" | /usr/local/bin/sshcad init'
