# systemd/sshcad.service
[Unit]
Description=SSH CA Certificate Signing Service
After=network.target

[Service]
Type=simple
User=sshcad
Group=sshcad

# Load environment variables from /etc/default/sshcad
EnvironmentFile=/etc/default/sshcad

# Systemd directory management
StateDirectory=sshcad
ConfigurationDirectory=sshcad

# Load CA private key from systemd credentials (TPM2-encrypted)
# Key stored in /etc/credstore.encrypted/ca-private-key (no .cred extension!)
# See README for secure key generation with zero disk writes
LoadCredentialEncrypted=ca-private-key:ca-private-key

ExecStartPre=/usr/bin/find $CREDENTIALS_DIRECTORY
ExecStartPre=/usr/bin/ls -al $CREDENTIALS_DIRECTORY
ExecStart=/usr/local/bin/sshcad serve

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/sshcad

# Restart policy
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target
