# test/mkosi/trees/test/etc/systemd/system/sshca-serv.service.d/10-test-user.conf
[Service]
ExecStartPre=/bin/bash -c 'echo -e "testpass\\ntestpass" | /usr/local/bin/sshca-serv add-user testuser'
