# test/mkosi/trees/test/etc/systemd/system/sshca-serv-test-setup.service
[Unit]
Description=SSH CA Service Test Setup
Before=sshca-serv.service
DefaultDependencies=no
After=local-fs.target

[Service]
Type=oneshot
RemainAfterExit=yes

ExecStart=ssh-keygen -t ed25519 -f /tmp/test_ca -N "" -C "Test CA"
ExecStart=systemd-creds encrypt --with-key=auto --name=ca-private-key /tmp/test_ca /etc/credstore.encrypted/ca-private-key
ExecStart=openssl req -x509 -newkey rsa:2048 -keyout /etc/sshca-serv/key.pem -out /etc/sshca-serv/cert.pem -days 365 -nodes -subj "/CN=localhost"
ExecStart=chown sshca-serv:sshca-serv /etc/sshca-serv/cert.pem /etc/sshca-serv/key.pem

[Install]
WantedBy=multi-user.target
